# Setup_files/rc.local
#!/bin/sh -e

cd /home/sharukh/Desktop/
sh AEP.sh

exit 0
